# contrib/fundermaps-refresh-model.service
[Unit]
Description=FunderMaps Worker - Submit Refresh Models Job
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/fundermaps/config.env
ExecStart=/usr/bin/psql "postgresql://${FUNDERMAPS_DB_USER}:${FUNDERMAPS_DB_PASSWORD}@${FUNDERMAPS_DB_HOST}:${FUNDERMAPS_DB_PORT}/${FUNDERMAPS_DB_NAME}?sslmode=require" -c "INSERT INTO application.worker_jobs (job_type, status) VALUES ('refresh_models', 'pending')"
SyslogIdentifier=fundermaps-submit-refresh
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
