# contrib/fundermaps-process-worker-jobs.service
[Unit]
Description=FunderMaps Worker - Job Queue Processor
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Environment="IMAGE=fundermaps-worker:latest"
ExecStart=/usr/bin/podman run --rm --name fundermaps-worker --env-file /etc/fundermaps/config.env ${IMAGE} process_worker_jobs.py --poll-interval 60 --max-concurrent 1 --job-types load_dataset export_product send_mail process_mapset generate_pdf refresh_models cleanup_storage --log-simple
SyslogIdentifier=fundermaps-worker
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
